# fm.cron
5 11,15,22 * * * python FundMonitor.py >> cron.log
